# rc.local
#!/bin/bash

#cd /data/etl-trans && nohup pypy etl-trans.py -r -m 2>&1 >/dev/null &
nohup /data/etl-trans/start.sh 2>&1 >/dev/null &
